# old-redhat62/els-base/els-base.cron
#!/bin/bash

. /usr/lib/process-watcher

TestProcess inetd inetd
TestProcess syslogd syslog
